# .config/systemd/user/fusuma.service
[Unit]
Description=Fusuma Multitouch

[Service]
ExecStart=/usr/local/bin/fusuma

[Install]
WantedBy=default.target
